# scanexpress.service.template
[Unit]
Description=ScanExpress Flask Service
After=network.target

[Service]
Type=simple
User=scanexpress
Group=scanexpress
WorkingDirectory=/opt/scanexpress
# Optional override when config.ini is not in the default location:
# Environment="SCANEXPRESS_CONFIG_FILE=/etc/scanexpress.conf"
ExecStart=/opt/scanexpress/.venv/bin/python /opt/scanexpress/app.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
